# etc/systemd/system/docker-sabnzbd-update.service
[Unit]
Description=sabnzbd docker update
Wants=docker.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-sabnzbd.sh
