# etc/systemd/system/docker-sickchill-update.service
[Unit]
Description=sickchill docker update
Wants=docker.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-sickchill.sh
