# etc/systemd/system/docker-plex-update.service
[Unit]
Description=plex docker update
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-plex.sh
