# etc/systemd/system/netdata-update.service
[Unit]
Description=netdata update service

[Service]
Type=oneshot
ExecStart=/usr/src/netdata/netdata-updater.sh
