# etc/systemd/system/netdata-update.service
[Unit]
Description=netdata update
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/netdata/netdata-updater.sh

