# etc/systemd/system/netdata-update.timer
[Unit]
Description=netdata update timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=86400
