# etc/systemd/system/raid-speed-limits.timer
[Unit]
Description=raid-speed-limits

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
