# etc/systemd/system/raid-scrub-monitor.service
[Unit]
Description=raid-scrub-monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/raid-scrub-monitor.sh
