# etc/systemd/system/raid-speed-limits.service
[Unit]
Description=raid-speed-limits

[Service]
Type=oneshot
ExecStart=/usr/local/bin/raid-speed-limits.sh
