# etc/systemd/system/raid-speed-limits.timer
[Unit]
Description=raid-speed-limits

[Timer]
OnBootSec=20sec

[Install]
WantedBy=timers.target
