# etc/systemd/system/raid-scrub.service
[Unit]
Description=report

[Service]
Type=oneshot
ExecStart=/usr/local/bin/raid-scrub.sh
