# etc/systemd/system/report.timer
[Unit]
Description=report

[Timer]
OnBootSec=1min
OnCalendar=Mon 08:00

[Install]
WantedBy=timers.target
