# etc/systemd/system/report.service
[Unit]
Description=report

[Service]
Type=oneshot
ExecStart=/usr/local/bin/report.sh
